# nginx/docker/nginx-rotate.cron
0 4 * * * /usr/sbin/logrotate /app/logrotate.d/nginx-rotate.conf
